# systemd/jarvis-voice-client.service
[Unit]
Description=JARVIS Voice Client
After=network.target pipewire.service pulseaudio.service
Documentation=https://github.com/Ex8-ca/jarvis-voice-shell

[Service]
Type=simple
# Adjust these paths to match your install
Environment=HOME=/home/%u
WorkingDirectory=/home/%u/jarvis-voice-shell
ExecStart=/home/%u/jarvis-voice-shell/venv/bin/python /home/%u/jarvis-voice-shell/jarvis_voice_client.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
